# sown-core/files/etc/init.d/led
#!/bin/sh /etc/rc.common
# We do something different for LEDs, so we don't want the original LED script